# launchd/com.trade-analysis.executor-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trade-analysis.executor-poll</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/takueisaotome/PycharmProjects/trade-analysis/scripts/run_executor_poll.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday through Friday at 06:32 PT (1st poll: most OPG fills) -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>32</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>32</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>32</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>32</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>32</integer></dict>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/takueisaotome/PycharmProjects/trade-analysis</string>
    <key>StandardOutPath</key>
    <string>/Users/takueisaotome/PycharmProjects/trade-analysis/logs/launchd-executor-poll.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/takueisaotome/PycharmProjects/trade-analysis/logs/launchd-executor-poll.log</string>
</dict>
</plist>
